# etc/sv/bluetoothd/run
#!/bin/sh

exec /usr/lib/bluetooth/bluetoothd --nodetach
