# x11/service/xcape/run
#!/bin/sh
export DISPLAY=:0
export XAUTHORITY=/home/martin/.Xauthority
exec xcape -de 'Shift_L=KP_Add'
