# x11/service/xcape/run
#!/bin/sh
export DISPLAY=:0
export XAUTHORITY=/home/martin/.Xauthority
exec xcape -fe 'Shift_L=KP_Add' >/dev/null 2>&1
