# _old/sv/sshd/run
#!/bin/sh
ssh-keygen -A # Generate host keys if they don't already exist.
exec /usr/sbin/sshd -D
