# x11/service/xbanish/run
#!/bin/sh
export DISPLAY=:0
export XAUTHORITY=/home/martin/.Xauthority
export PATH=$PATH:/usr/local/bin
exec xbanish -i mod4 >/dev/null 2>&1
