# x11/service/xbanish/run
#!/bin/sh
export DISPLAY=:0
export XAUTHORITY=/home/martin/.config/x11/authority
exec xbanish -i mod4 >/dev/null 2>&1
