# etc/sv/runsvdir-martin/run
#!/bin/sh
exec 2>&1
exec chpst -umartin runsvdir /home/martin/.config/service
