# service/xbanish/run
#!/bin/sh
export DISPLAY=:0
export XAUTHORITY=/home/martin/.Xauthority
exec xbanish -i mod4
